# vixeditor-webhook.service
[Unit]
Description=VixEditor Webhook Service (Production)
After=network.target

[Service]
Type=notify
User=root
WorkingDirectory=/home/ubuntu/vixeditor
Environment="PATH=/home/ubuntu/vixeditor/venv/bin"

# Use Gunicorn for production (not Flask dev server)
ExecStart=/home/ubuntu/vixeditor/venv/bin/gunicorn -c gunicorn_config.py deploy_webhook:app

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/ubuntu/vixeditor/logs/webhook-service.log
StandardError=append:/home/ubuntu/vixeditor/logs/webhook-service-error.log

[Install]
WantedBy=multi-user.target
