# vixeditor-webhook.service
[Unit]
Description=VixEditor Webhook Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/CatVideo
Environment="PATH=/path/to/CatVideo/venv/bin"
ExecStart=/path/to/CatVideo/venv/bin/python deploy_webhook.py
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/path/to/CatVideo/logs/webhook-service.log
StandardError=append:/path/to/CatVideo/logs/webhook-service-error.log

[Install]
WantedBy=multi-user.target
